# numen/numen-0.7/service/runit/run
#!/bin/sh
exec 2>&1
[ -r ./conf ] && . ./conf
: ${NUMENUSER:=_numen}
: ${NUMENGROUP:=_numen}
if [ -z "$NUMENHOME" ]; then
	mkdir -p /var/lib/numen
	chown "$NUMENUSER:$NUMENGROUP" /var/lib/numen
	export XDG_STATE_HOME=/var/lib
fi
exec chpst -u "$NUMENUSER:$NUMENGROUP:input:audio" numen ${OPTS}
