# ap/slint-user-settings/slint-user-settings/etc/skel/.config/autostart/startorca.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=orca
Exec=/usr/bin/orca
Hidden=true
